# netbsd/etc/rc.d/mobileip6
#!/bin/sh
# $Id: mobileip6,v 1.7 2005/05/25 10:07:30 keiichi Exp $

# REQUIRE: network
# PROVIDE: mobileip6
# BEFORE: LOGIN

. /etc/rc.subr

name="mobileip6"
start_cmd="mobileip6_start"
stop_cmd="mobileip6_stop"


eval ipv6_mobile_enable=${ipv6_mobile_enable:=NO}
eval ipv6_mobile_config_dir=${ipv6_mobile_config_dir:=/usr/local/v6/etc/mobileip6}

eval ipv6_mobile_mobile_node_interfaces=${ipv6_mobile_mobile_node_interfaces:=mip0}
eval ipv6_mobile_mnd_flags=${ipv6_mobile_mnd_flags:='"-n -c /usr/local/v6/etc/mnd.conf"'}
eval ipv6_mobile_mrd_flags=${ipv6_mobile_mrd_flags:='"-n -c /usr/local/v6/etc/mrd.conf"'}

eval ipv6_mobile_mdd_program=${ipv6_mobile_mdd_program:=babymdd}

eval ipv6_mobile_home_agent_mode=${ipv6_mobile_home_agent_mode:=mobile_ipv6}
eval ipv6_mobile_had_flags=${ipv6_mobile_had_flags:='"-n -c /usr/local/v6/etc/had.conf"'}

eval ipv6_mobile_cnd_flags=${ipv6_mobile_cnd_flags:='"-n -c /usr/local/v6/etc/cnd.conf"'}

mobileip6_start() 
{
	echo -n 'Doing Mobile IPv6 network setup:'

	case ${ipv6_mobile_debug_enable} in
	[Yy][Ee][Ss])
		sysctl -w net.inet6.mip6.debug=1
		;;
	[Nn][Oo])
		sysctl -w net.inet6.mip6.debug=0
		;;
	esac

	case ${ipv6_mobile_security_enable} in
	[Yy][Ee][Ss])
		sysctl -w net.inet6.mip6.use_ipsec=1
		;;
	[Nn][Oo])
		sysctl -w net.inet6.mip6.use_ipsec=0
		;;
	esac

	case ${ipv6_mobile_nodetype} in
	[Mm][Oo][Bb][Ii][Ll][Ee]_[Nn][Oo][Dd][Ee])
		# setup SPD/SA
		if [ `sysctl -n net.inet6.mip6.use_ipsec` = '1' ]; then
			for node_dir in ${ipv6_mobile_config_dir}/*
			do
				if [ ! -e ${node_dir}/add ]; then
					continue;
				fi
				setkey -f ${node_dir}/add
				setkey -f ${node_dir}/spdadd_mobile_node
			done
		fi

		echo -n ' mnd'
		echo -n " ${ipv6_mobile_mdd_program}"
		for mip_interface in ${ipv6_mobile_mobile_node_interfaces}
		do
			mnd ${ipv6_mobile_mnd_flags} ${mip_interface}
			${ipv6_mobile_mdd_program} -h ${mip_interface}
		done
		echo -n ' rtsold'
		sysctl -w net.inet6.ip6.accept_rtadv=1
		rtsold -a -m
		echo -n ' cnd'
		cnd ${ipv6_mobile_cnd_flags}
		;;
	[Mm][Oo][Bb][Ii][Ll][Ee]_[Rr][Oo][Uu][Tt][Ee][Rr])
		# setup SPD/SA
		if [ `sysctl -n net.inet6.mip6.use_ipsec` = '1' ]; then
			for node_dir in ${ipv6_mobile_config_dir}/*
			do
				if [ ! -e ${node_dir}/add ]; then
					continue;
				fi
				setkey -f ${node_dir}/add
				setkey -f ${node_dir}/spdadd_mobile_node
			done
		fi

		echo -n ' mrd'
		echo -n " ${ipv6_mobile_mdd_program}"
		echo -n ' nemonetd'
		nemonetd_flags=${ipv6_mobile_nemonetd_flags:='-n -m -c /usr/local/v6/etc/mrd.conf'}
		for mip_interface in ${ipv6_mobile_mobile_node_interfaces}
		do
			mrd ${ipv6_mobile_mrd_flags} ${mip_interface}
			${ipv6_mobile_mdd_program} -h ${mip_interface}
			nemonetd ${nemonetd_flags} ${mip_interface}
		done
		echo -n ' rtsold'
		sysctl -w net.inet6.ip6.accept_forwarding=1
		sysctl -w net.inet6.ip6.accept_rtadv=1
		rtsold -a -m 
		echo -n ' cnd'
		cnd ${ipv6_mobile_cnd_flags}
		;;
	[Hh][Oo][Mm][Ee]_[Aa][Gg][Ee][Nn][Tt])
		case ${ipv6_mobile_home_interface} in
		'')
			;;
		*)
			# setup SPD/SA
			if [ `sysctl -n net.inet6.mip6.use_ipsec` = '1' ]; then
				for node_dir in ${ipv6_mobile_config_dir}/*
				do
					if [ ! -e ${node_dir}/add ]; then
						continue;
					fi
					setkey -f ${node_dir}/add
					setkey -f ${node_dir}/spdadd_home_agent
				done
			fi

			home_prefixes=`mobileip6_getprefixes ${ipv6_mobile_home_interface}`
			for home_prefix in ${home_prefixes}
			do
				# XXX describe the case when plen != 64
				ifconfig ${ipv6_mobile_home_interface} inet6 ${home_prefix}:fdff:ffff:ffff:fffe prefixlen 64 anycast alias
			done

			echo -n ' had'
			had ${ipv6_mobile_had_flags} ${ipv6_mobile_home_interface}
			if [ ${ipv6_mobile_home_agent_mode} = "nemo" ]; then
				echo -n ' nemonetd'
				nemonetd_flags=${ipv6_mobile_nemonetd_flags:='-n -h -c /usr/local/v6/etc/had.conf'}
				nemonetd ${nemonetd_flags} ${ipv6_mobile_home_interface}
			fi
#			echo -n ' cnd'
#			cnd ${ipv6_mobile_cnd_flags}
			echo -n ' rtadvd'
			rtadvd -m ${ipv6_mobile_home_interface}
			# take in Router Advertisement
			rtsol ${ipv6_mobile_home_interface} > /dev/null 2>&1
			;;
		esac
		;;
	*)
		echo -n ' cnd'
		cnd ${ipv6_mobile_cnd_flags}
		;;
	esac

	echo ''
}

mobileip6_stop() {

}

mobileip6_getprefixes() {
	ndp -p 2>/dev/null | while read prefix interface rest; do
		if [ X$interface = X"if="$1 ]; then
			case ${prefix} in
			fe80:*)
				;;
			fec0:*)
				;;
			[0-9a-f][0-9a-f][0-9a-f][0-9a-f]:*)
				prefix64=`expr "${prefix}" : '\(.*\)::/64'`
				echo ${prefix64}
			esac
		fi
	done
}

load_rc_config $name
run_rc_command "$1"
